# images/tt-rss/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m  h dom mon dow command
*/30 * *   *   *   /usr/bin/php /var/www/tt-rss/update.php --feeds --quiet
